# dotfiles/config/.config/systemd/user/check-claude-status.timer
[Unit]
Description=Check Claude API status timer
BindsTo=graphical-session.target
After=graphical-session.target

[Timer]
OnStartupSec=300
OnUnitActiveSec=300

[Install]
WantedBy=graphical-session.target
