# dotfiles/config/.config/systemd/user/check-claude-status.service
[Unit]
Description=Check Claude API status
BindsTo=graphical-session.target
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/environment/bin/check_claude_status
