# dotfiles/config/.config/systemd/user/x-session.target
[Unit]
Description=X11 graphical session
Documentation=man:systemd.special(7)
